# iam_user_sync.timer
[Unit]
Description=Run iam_user_sync.service every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
